# etc/services.d/haproxy/run
#!/usr/bin/with-contenv sh

# fixme check if empty

exec s6-setuidgid haproxy -f /etc/haproxy/haproxy.conf;
