# etc/services.d/haproxy/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid balancer haproxy -f /etc/haproxy/haproxy.conf;
